# systemd/av-voice.service
[Unit]
Description=Runs /home/pi/av_dev/voice.py
  
[Service]
ExecStart=/home/pi/av_dev/voice.py

[Install]
WantedBy=multi-user.target
